# sound-station/birdnet-export.timer
[Unit]
Description=Run the BirdNET-Go S3 exporter every minute

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
